# black-hole/cron.daily
#!/bin/bash

cd ~/funny/

git pull
if [ $? -ne 0 ]; then exit 1; fi

git commit -m 'log sync' black-hole/blackHole.log
if [ $? -ne 0 ]; then exit 1; fi

git push
